# src/port-health-checker-.service.in
[Unit]
Description=Check connectivity to the network port @PORT@ on @ADDR@
Documentation=man:nc(1)
OnFailure=notification-mail@%n.service

[Service]
Type=oneshot
# TODO:
#EnvironmentFile=/etc/port-health-checker/%i.conf
ExecStart=/usr/bin/nc --nodns -z $OPTIONS @ADDR@ @PORT@
Wants=nginx.service
After=nginx.service

# Hardening options:
PrivateDevices=true
PrivateTmp=true
RestrictAddressFamilies=AF_INET AF_INET6
# .. note:: These are not available in RHEL =< 8.
#MemoryDenyWriteExecute=yes
#PrivateMounts=true
#RestrictRealtime=yes
#RestrictSUIDSGID=yes
#ProtectControlGroups=true
