# src/port-health-checker-.timer.in
[Unit]
Description=Port health chekcer

[Timer]
OnCalendar=hourly
Persistent=true
AccuracySec=24h

[Install]
WantedBy=timers.target
