# src/port-health-checker@.timer
[Unit]
Description=Port health chekcer

[Timer]
OnCalendar=hourly
AccuracySec=24h

[Install]
WantedBy=timers.target
